# config/systemd/user/geoclue-agent.service
[Unit]
Description=Geoclue2 User Agent

[Service]
ExecStart=/usr/lib/geoclue-2.0/demos/agent

[Install]
WantedBy=default.target
